# cloudflareddns.timer
[Unit]
Description=Run Cloudflare DDNS Updater

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
